# services/monitor.service
[Unit]
Description=RAID Monitor Service

[Service]
ExecStart=/home/boike/raidmonitor/monitor.sh